# src/configs/profile.d/owb.sh
export PATH=$PATH:/usr/local/owb/backend/usr/bin
